# ledify.service
[Unit]
Description=Ledify
After=pilight.service

[Service]
Restart=always
ExecStart=/usr/local/bin/ledify

[Install]
WantedBy=multi-user.target
