# ledify.service
[Unit]
Description=Ledify

[Service]
ExecStart=/usr/local/bin/ledify

[Install]
WantedBy=multi-user.target
